# supervisor/bing.conf
[program:postgres]
command = uwsgi --ini /var/www/postgres/bing_uwsgi.ini
stopsignal=QUIT
autostart=true
autorestart=true
redirect_stderr=true
